# bin/mdroid-core.service
[Unit]
Description=MDroid Core
After=influxd.service
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=2
User=root
ExecStart=/home/pi/go/bin/GoQMW --settings-file /home/pi/MDroid/config.json &> /home/pi/MDroid/logs/goqmw/`date +\%Y-\%m-\%d_\%H:\%M:\%S`-goqmw.log

[Install]
WantedBy=multi-user.target